# plexamp-lrc.service
[Unit]
Description=Plexamp LRC++ Lyrics Manager
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/plexamp-lrc-plusplus
EnvironmentFile=/home/pi/plexamp-lrc-plusplus/.env
ExecStart=/home/pi/plexamp-lrc-plusplus/venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
